# start-ipkvm.socket
[Unit]
Description=OpenBMC ipKVM daemon socket

[Socket]
ListenStream=5900

[Install]
WantedBy=sockets.target
